# com.whisper.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.whisper.service</string>
    <key>Program</key>
    <string>/Volumes/workplace/tools/whisper/.venv/bin/python3</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/workplace/tools/whisper/.venv/bin/python3</string>
        <string>/Volumes/workplace/tools/whisper/whisper_service.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/jainnam/Library/Logs/com.whisper.service.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jainnam/Library/Logs/com.whisper.service.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Volumes/workplace/tools/whisper</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/jainnam/.local/share/mise/installs/python/3.13.6/bin</string>
        <key>PYTHONPATH</key>
        <string>/Volumes/workplace/tools/whisper</string>
    </dict>
</dict>
</plist> 