# com.whisper.hostkey.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.whisper.hostkey</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/workplace/tools/whisper/run.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Volumes/workplace/tools/whisper</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PICOVOICE_ACCESS_KEY</key>
        <string>lk++IHEpUel5qLDl6dc4e2qR12RqlKoMNzILpflCnLYVuTba4t3v0w==</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/tmp/whisper-hostkey.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/whisper-hostkey.error.log</string>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>Nice</key>
    <integer>-10</integer>
</dict>
</plist>