# whisper_settings.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.whisper-dictation</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/Workspace/whisper-dictation/startup.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Volumes/Workspace/whisper-dictation/whisper-dictation.log</string>
    <key>StandardErrorPath</key>
    <string>/Volumes/Workspace/whisper-dictation/whisper-dictation.err</string>
    <key>WorkingDirectory</key>
    <string>/Volumes/Workspace/whisper-dictation</string>
</dict>
</plist> 